# root/etc/services.d/deez/run
#!/usr/bin/with-contenv bash

echo "Starting Deez..."
exec \
	s6-setuidgid abc node \
	"/app/deez/app/app.js"
